# snips/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -e
stdout_logfile=/data/debian-base/log/%(program_name)s.log
stderr_logfile=/data/debian-base/log/%(program_name)s.log
autorestart=true
